# templates/supervisord.conf.erb
[supervisord]
nodaemon=true

[program:dnsmasq]
command=/usr/sbin/dnsmasq
